# examples/uwsgi.service
[Unit]
Description = uWSGI
After = syslog.target

[Service]
ExecStart = /home/pi/temperaturess/bin/uwsgi --ini /home/pi/Temperaturess/uwsgi.ini
Restart=always
KillSignal=SIGQUIT
Type=notify
StandardError=syslog
NotifyAccess=all

[Install]
WantedBy=multi-user.target